# binaire/agent/minden-agent.service
[Unit]
Description=Minden Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=SERVER_URL=https://your.minden.app
Environment=TOKEN=REPLACE_WITH_AGENT_TOKEN
ExecStart=/opt/minden-agent/minden-agent
Restart=always
RestartSec=5
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target

